# tempura.service
[Unit]
Description=tempura
 
[Service]
ExecStart=/usr/bin/tempura
 
[Install]
WantedBy=multi-user.target
